# systemd-units/resticbackup.timer
[Unit]
Description=Timer for Restic Backup script

[Timer]
OnCalendar=*-*-* 00:00:00
OnStartupSec=60m
RandomizedDelaySec=120m
Persistent=true

[Install]
WantedBy=timers.target